# ddind/loopdevice.service
[Unit]
Description=Setup Loop Devices

[Service]
Type=oneshot
ExecStart=/bin/sh /init.sh

[Install]
WantedBy=multi-user.target
